# autossh.service
[Unit]
Description=AutoSSH to SSH Gateway
After=network-online.target

[Service]
user=autossh
# This file goes in /etc/systemd/system/
# -p [PORT]
# -l [user]
# -M 0 --> no monitoring
# -N Not interactive
# LOCALPORT:IP_ON_EXAMPLE:PORT_ON_EXAMPLE
ExecStart=/usr/bin/autossh -M 0 -N -o "ServerAliveInterval 5" -o "ServerAliveCountMax 3" -p 22 -l **REMOTEUSER** **REMOTEHOST** -R 2201:127.0.0.1:22 -i /home/autossh/.ssh/id_rsa

RestartSec=5
Restart=always

[Install]
WantedBy=multi-user.target
